# etc/systemd/system/podman-prune.service
[Unit]
Description=Podman prune

[Service]
Type=oneshot
User=root
ExecStart=podman system prune -af --build

[Install]
WantedBy=multi-user.target
